# supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=/usr/sbin/mysqld
user=mysql
autorestart=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true
